# conf/compotes.timer
[Unit]
Description=Schedule compote reminder

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
